# config/services/udpgw-7100.service
[Unit]
Description=BadVPN UDP Gateway for Port 7100
After=network.target

[Service]
ExecStart=/usr/bin/badvpn-udpgw --listen-addr 127.0.0.1:7100 --max-clients 500
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
